# com.wctray.sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.wctray.sync</string>

    <key>ProgramArguments</key>
    <array>
        <string>SCRIPT_DIR/.venv/bin/python</string>
        <string>SCRIPT_DIR/sync_schedule.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>SCRIPT_DIR</string>

    <key>StartCalendarInterval</key>
    <array>
        <dict>
            <key>Hour</key>
            <integer>6</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
        <dict>
            <key>Hour</key>
            <integer>23</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
    </array>

    <key>StandardOutPath</key>
    <string>SCRIPT_DIR/sync.log</string>

    <key>StandardErrorPath</key>
    <string>SCRIPT_DIR/sync.log</string>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
